# hello_world2/rootfs/etc/services.d/hello_world2/run
#!/usr/bin/python3
 
print("Hello World2b from python run2.py!")