# hello_world2/rootfs/etc/services.d/hello_world2/run
#!/usr/bin/python3

import os
 
print("Hello World2b from python run2.py!")
os.system("ping -c 10 192.168.0.230")