# hello_world2/rootfs/etc/services.d/hello_world2/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================


declare message
message=$(bashio::config 'message')
bashio::log.info "${message:="run5 <Hello_World2>"}"

#cat /etc/crontabs/root
#crontab -l
python3 /opt/run2.py
