# hello_world2/rootfs/etc/services.d/hello_world2/run
#!/usr/bin/with-contenv bashio

#echo Hello world2c... pinging!
#ping -c 1 192.168.0.230

python3 /opt/run2.py

#ping -c 1 192.168.0.230
#echo completed!
