# discover_devices/rootfs/etc/services.d/hello_world2/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================


declare apikey
apikey=$(bashio::config 'apikey')
bashio::log.info "APIKey: ${apikey:=""}"

#Execute every 10 seconds
while true
do
    python3 discover_device.py
    sleep $(( 2*2 ))
done
